# rootfs/etc/systemd/system/kazeta-profile-loader.service
[Unit]
Description=Load custom InputPlumber profile for Kazeta+
# We need to wait for both InputPlumber AND the gamer user session to be ready.
Requires=inputplumber.service user@1000.service
After=inputplumber.service user@1000.service

[Service]
Type=oneshot
# The small delay is still a good idea to prevent any race conditions
ExecStartPre=/bin/sleep 2
# This is the exact busctl command you provided
ExecStart=/usr/bin/busctl call org.shadowblip.InputPlumber /org/shadowblip/InputPlumber/CompositeDevice0 org.shadowblip.Input.CompositeDevice LoadProfilePath "s" /usr/share/inputplumber/profiles/steam-deck.yaml

[Install]
WantedBy=multi-user.target
